# shutdown-after-sleep-installer.service
[Unit]
Description=Reinstall shutdown-after-sleep hook (steam-deck-smart-sleep) after SteamOS update
After=local-fs.target

[Service]
Type=oneshot
ExecStart=__INSTALL_SCRIPT__
RemainAfterExit=yes
